# lib/lb2.service
[Unit]
Description=XUI LoadBalance - Motor de entrega de streaming
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
ExecStart=/home/xui/loadbalance/lb2 -config /home/xui/loadbalance/config.json
Restart=always
RestartSec=3

# Cada espectador consome um file descriptor; o padrao de 1024 estoura de imediato.
LimitNOFILE=1000000

# O processo so precisa ler a propria config, os segmentos do XUI e o MySQL.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

StandardOutput=append:/var/log/lb2.log
StandardError=append:/var/log/lb2.log

[Install]
WantedBy=multi-user.target
